# .etc-backup-config-files/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin rmw --noclear %I $TERM
